# archive/configs/supervisor-single.conf
[supervisord]
nodaemon=true

[program:safework]
command=/app/start-single.sh
autostart=true
autorestart=true
stdout_logfile=/app/logs/app.log
stderr_logfile=/app/logs/app_error.log